# config/systemd/atlas-daily-maintenance.service.in
[Unit]
Description=Atlas daily complete maintenance loop
StartLimitIntervalSec=3h
StartLimitBurst=3

[Service]
Type=exec
WorkingDirectory=@PROJECT_ROOT@
ExecStart=@PROJECT_ROOT@/scripts/run_daily_maintenance.sh
Restart=no
RestartForceExitStatus=75
RestartSec=20min
RuntimeMaxSec=14h
TimeoutStopSec=12min
